# ansible/roles/resource_usage/templates/resource-usage.service.j2
[Unit]
Description=Resource Usage

[Service]
User=safe
ExecStart=/home/safe/resource-usage.sh
NetworkAccess=yes

[Install]
WantedBy=multi-user.target
